# install-init/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/tty1message
StandardInput=tty
StandardOutput=tty
